# sys-auth/validity-sensor-daemon/files/validity-sensor
#!/sbin/openrc-run
# Copyright 1999-2013 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

daemon_binary="vcsFPService"

depend() {
    after modules
}

start() {
	ebegin "Starting Validity Fingerprint sensor daemon"
	rm -f ${pidfile}
	start-stop-daemon --start --quiet --exec /usr/bin/${daemon_binary}
	eend $?
}

stop() {
	ebegin "Stopping Validity Fingerprint sensor daemon"
    pkill -TERM ${daemon_binary}
	eend $?
}
